# systemd/rtorrent-cleanup.timer
[Unit]
Description=Run rtorrent-cleanup periodically

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
